# automation/gitdigest.service
[Unit]
Description=GitDigest daily report
After=network.target

[Service]
Type=oneshot
User=your_username
WorkingDirectory=/path/to/GitDigest
Environment=GITHUB_TOKEN=your_token_here
ExecStart=/usr/bin/python3 run_digest.py
StandardOutput=append:/path/to/GitDigest/output/service.log
StandardError=append:/path/to/GitDigest/output/service.log

[Install]
WantedBy=multi-user.target
